# extras/systemd/foreman.service
[Unit]
Description=Foreman
Documentation=https://theforeman.org
After=network.target remote-fs.target nss-lookup.target
Requires=foreman.socket

[Service]
Type=simple
User=foreman
TimeoutSec=300
WorkingDirectory=/usr/share/foreman
ExecStart=/usr/share/foreman/bin/rails server --environment $FOREMAN_ENV
Environment=FOREMAN_ENV=production FOREMAN_BIND=tcp://0.0.0.0:3000
Environment=MALLOC_ARENA_MAX=2

SyslogIdentifier=foreman

[Install]
WantedBy=multi-user.target
